# cron.d
* * * * * /home/goez/.cla/eggoez >/dev/null 2>&1
